# ROS/src/BUD-E/hearbo_2dnav/hearbo_cart_ts/hearbo_cart_ts_ctrl/build/devel/etc/catkin/profile.d/05.catkin-test-results.sh
# generated from catkin/cmake/env-hooks/05.catkin-test-results.sh.develspace.in

export CATKIN_TEST_RESULTS_DIR="/home/praveen/Documents/ROS/tmp/hearbo_2dnav/hearbo_cart_ts/hearbo_cart_ts_ctrl/build/test_results"
export ROS_TEST_RESULTS_DIR="$CATKIN_TEST_RESULTS_DIR"
